# deploy/systemd/okr-stack-centos.service.example
[Unit]
Description=OKR stack with embedded kkFileView runtime (CentOS)
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/srv/okr
EnvironmentFile=-/etc/sysconfig/okr-stack
ExecStart=/srv/okr/scripts/linux/start-all.sh
ExecStop=/srv/okr/scripts/linux/stop-all.sh
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
